# ebyt.service
[Unit]
Description=EBYT — Every Byte You Track (activity tracker)
After=graphical-session.target

[Service]
ExecStart=/usr/bin/ebyt daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
